# deploy-cn/systemd/pms-api.service
[Unit]
Description=PMS API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/pms/api
ExecStart=/usr/bin/dotnet /opt/pms/api/PMS.API.dll
Environment=ASPNETCORE_URLS=http://127.0.0.1:5111
Environment=ASPNETCORE_ENVIRONMENT=Production
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
